# ansible/templates/speechtractor.service.j2
# speechtractor.service

[Unit] 
Description=Speechtractor server

[Service] 
Type=simple 
ExecStart=/usr/services/speechtractor_service.sh
